# deploy/attendance.service
[Unit]
Description=Attendance Audit Service
After=network.target postgresql.service

[Service]
Type=simple
User=attendance
WorkingDirectory=/opt/attendance
ExecStart=/usr/bin/java -jar /opt/attendance/target/attendance-audit-0.0.1-SNAPSHOT.jar --spring.profiles.active=prod
Restart=always
RestartSec=10
StandardOutput=append:/opt/attendance/logs/attendance.out
StandardError=append:/opt/attendance/logs/attendance.err

[Install]
WantedBy=multi-user.target
